# wbc-supervisor.conf
[program:wbc-macbre-net]
autostart = true
command = make start
directory = /home/macbre/github/wbc.macbre.net
user = macbre
